# systemd/system/null-routes.timer
[Unit]
Description=Run null-routes every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=10s

[Install]
WantedBy=timers.target
